# linux/.config/systemd/user/suspend@.timer
[Unit]
Description=Suspend

[Timer]
Unit=suspend.service
OnCalendar=*-*-* %I

[Install]
WantedBy=timers.target
